# packages/ue-operator_skeleton/etc/systemd/system/ue-operator.service
[Unit]
Description=ue-operator Service
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
User=ubilab
ExecStart=/usr/bin/ue-operator start

[Install]
WantedBy=multi-user.target
