# packages/ue-operator/etc/systemd/system/ue-operator.service
[Unit]
Description=ue-operator Service

[Service]
Type=simple
ExecStart=/usr/bin/ue-operator start

[Install]
WantedBy=multi-user.target
